# pacmauto.service
[Unit]
Description=Pacmauto Service

[Service]
Type=simple
# ExecStart=/the/path/to/the/pacmauto/script

[Install]
WantedBy=multi-user.target
